# var/www/html/alpine/config/a0:36:9f:e7:6a:e2/etc/init.d/acontroll
#!/sbin/openrc-run

description="Sets some parameters of the machine."

depend() {
        keyword -prefix -lxc -docker
}

start() {
        #if [ -s /etc/hostname ] ; then
        #        opts="-F /etc/hostname"
        #else
        #        opts="${hostname:-localhost}"
        #fi
        ebegin "Setting parameters"
	sleep 10;
        #hostname $opts
        eend $?
}
